# systemd/lbcd.service.in
# systemd configuration for lbcd.  -*- conf -*-

[Unit]
Description=responder for load balancing
Documentation=man:lbcd(8)
Documentation=http://www.eyrie.org/~eagle/software/lbcd/
Requires=lbcd.socket

[Service]
Type=simple
NotifyAccess=main
Restart=always
User=lbcd
Group=lbcd

# Support for /etc/default/lbcd is here only for backward compatibility with
# the lbcd init script shipped with Debian.  It and the $DAEMON_OPTS reference
# can be dropped if this backward compatibility is not needed.
EnvironmentFile=-/etc/default/lbcd
ExecStart=@sbindir@/lbcd -f -l $DAEMON_OPTS

[Install]
Also=lbcd.socket
WantedBy=multi-user.target
